# ERICnrcloudran/src/main/resources/crontabs/cronfile
0 2 * * * /opt/ericsson/nrcloudransystemtopology/cron.sh 2>/dev/null
